# artifacts/rc.local
#!/bin/bash -e
sudo -i -u ubuntu /home/ubuntu/start_test_harness.sh
